# system/efistub-ucode-update.path
[Unit]
Description=Copy Intel microcode to UEFISYS Partition

[Path]
PathChanged=/boot/intel-ucode.img

[Install]
WantedBy=multi-user.target
