# system/refind-update.service
[Unit]
Description=Install rEFInd update

[Service]
Type=oneshot
ExecStart=/usr/bin/refind-install
