# system/refind-update.path
[Unit]
Description=Install rEFInd update

[Path]
PathChanged=/usr/share/refind

[Install]
WantedBy=multi-user.target
